# ctf-levels.conf
description "creation and start of docker containers for stripe-ctf2"
version "1.0"
start on started dns4docker-upstart
chdir /opt/levels
script
exec /usr/local/bin/docker-compose up
end script

