# dns/dns4docker-upstart.conf
description "dynamic updates for docker containers in bind"
version "1.0"
start on filesystem
script
exec /usr/bin/dns4docker.pl
end script

